# roles/osmetrics/templates/osmetrics.service.j2
[Unit]
Description=osmetrics
After=network.target

[Service]
ExecStart={{ osmetrics_install_dir }}/osmetrics
Type=simple
PIDFile=/var/run/osmetrics.pid

[Install]
WantedBy=default.target
